# crontab
# START CRON JOB
* * * * * ./monitor.py --fqdn smtp.koeroo.net --port 25 --protocol SMTP > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB